# server-scripts/recordings-api.service
[Unit]
Description=Recordings Delete API
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/recordings-api.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target
